# crontab
17,37,57 * * * * /afs/nd.edu/user5/dlopes/Public/therm 129.74.152.75